# resources/Synology/yamaha.conf
description "Yamaha controller"
start on syno.share.ready
stop on runlevel [06]
setuid http
exec node /volume1/web/yamaha/app.js
